# deploy/ai-robot-dashboard-snapshot.timer
[Unit]
Description=Publish sanitized AI Robot dashboard snapshot every 15 minutes

[Timer]
OnBootSec=2min
OnCalendar=*:0/15
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
