# deploy/ai-robot-paper.service
[Unit]
Description=AI Robot OKX BTC/ETH paper swap trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/typbing/ai_robot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/typbing/.config/ai_robot/env
ExecStart=/usr/bin/python3 -u -m ai_robot.runner loop --config config.paper.json
Restart=always
RestartSec=15

[Install]
WantedBy=default.target
